# supervisor.conf
[supervisord]
nodaemon=true

[program:arrplat-server]
command=python /app/main.py

[program:arrplat-admin]
command=yarn start
